# AllStar-build/common/asterisk.timer
#/etc/systemd/system/asterisk.timer
[Timer]
OnStartupSec=10

[Install]
WantedBy=multi-user.target
